# config/defaults/service/php-fpm/run
#!/bin/bash
exec 2>&1
echo "starting php-fpm"
exec /usr/sbin/php-fpm7 -F
