# config/defaults/service/nginx/run
#!/bin/bash
exec 2>&1
echo "starting nginx"
[[ -f ${SHARED}/sock/nginx.sock ]] && rm ${SHARED}/sock/nginx.sock
exec /usr/sbin/nginx
