# config/defaults/service/crond/run
#!/bin/bash
exec 2>&1
echo "starting cron"
exec /usr/sbin/crond -f
